# base-system/etc/moxa-configs/UC-8220-LX/udev/10-gps.rules
SUBSYSTEM=="tty", ATTRS{idVendor}=="1546", ATTRS{idProduct}=="01a8", ACTION=="add", RUN+="/bin/bash -c '/bin/stty -F /dev/%k -echo'" 
