# base-system/09-watchdog.rules
KERNEL=="rtc0", RUN+="/sbin/modprobe ds1374_wdt"
